# cnc-sniffer.conf
# Pcap uploader

description     "cnc sniffer"

start on runlevel [2345]
stop on runlevel [016]

respawn

umask 007
chdir /opt/cnc-sniffer

script

        python main.py

end script



